# ossec-rules/windows/process_creation/win_win10_sched_task_0day.xml
<rule id="262220" level="14">
	<if_group>sysmon_event1</if_group>
	<field name="win.eventdata.Image">\\\\schtasks.exe$</field>
	<field name="win.eventdata.CommandLine">/change\.*/TN\.*/RU\.*/RP</field>
	<description>ATT&CK T1053: Windows 10 Scheduled Task SandboxEscaper 0-day</description>
	<info type="text">Detects Task Scheduler .job import arbitrary DACL write\par </info>
	<info type="text">Falsepositives: Unknown. </info>
	<info type="text">Sigma UUID: 931b6802-d6a6-4267-9ffa-526f57f22aaf </info>
	<info type="link">https://github.com/SandboxEscaper/polarbearrepo/tree/master/bearlpe </info>
	<group>attack.privilege_escalation,attack.execution,attack.t1053,car.2013-08-001,MITRE</group>
</rule>

